# rootfilesystem/etc/workers/worker.ini
[program:ladang-worker]
process_name=%(program_name)s
command=php /home/ladang/app/artisan horizon
user=ladang
autostart=true
autorestart=true
stopasgroup=true
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=3600
