# etc/crontabs/ladang
*       *       *       *       *       php /home/ladang/app/artisan schedule:run >> /dev/null 2>&1
0       *       *       *       *       php /home/ladang/app/artisan queue:restart >> /dev/null 2>&1
