# rootfilesystem/etc/supervisor.d/swoole.ini
[program:swoole]
command=php -d variables_order=EGPCS /home/ladang/app/artisan octane:start --host=0.0.0.0 --workers=%(ENV_OCTANE_WORKER)s --task-workers=%(ENV_OCTANE_TASK_WORKER)s --max-requests=%(ENV_OCTANE_MAX_REQUESTS)s
user=ladang
autostart=true
autorestart=true
priority=5
startretries=3
stopwaitsecs=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
