# rootfilesystem/etc/crontabs/ladang
*       *       *       *       *       php /home/ladang/app/artisan schedule:run >> /dev/null 2>&1
